# etc/systemd/system/ddns-onboot@.timer
[Unit]
Description=Activate DDNS after boot

[Timer]
OnActiveSec=45
RemainAfterElapse=no

[Install]
WantedBy=named.service

